# contrib/systemd/app.openprescribing.web.socket
[Socket]
ListenStream=/webapps/openprescribing/run/gunicorn.sock

[Install]
WantedBy=sockets.target
