# contrib/systemd/app.openprescribing_staging.web.socket
[Socket]
ListenStream=/webapps/openprescribing_staging/run/gunicorn.sock

[Install]
WantedBy=sockets.target
